# supervisor/supervisor/gunicorn.conf
[program:gunicorn]
command=gunicorn -c server.conf app:app
directory=/var/www/app
user=nobody
autostart=true
autorestart=true
redirect_stderr=true